# DEBIAN/postinst
#!/bin/bash
set -e

# Verificar se o usuário que executa o script é root
if [ "$(id -u)" -ne 0 ]; then
    echo "O script postinst deve ser executado como root."
    exit 1
fi

# Obter o nome do usuário do diretório HOME padrão
USER=$(logname)
if [ -z "$USER" ]; then
    echo "Não foi possível determinar o usuário atual."
    exit 1
fi

# Configurar o ambiente do usuário para o systemctl --user
XDG_RUNTIME_DIR="/run/user/$(id -u "$USER")"
export XDG_RUNTIME_DIR

# Executar os comandos no contexto do usuário
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user daemon-reload"
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user enable md-tracker-search-provider.service"
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user start md-tracker-search-provider.service"