# DEBIAN/prerm
#!/bin/bash
set -e

# Verificar se o XDG_RUNTIME_DIR está configurado corretamente
if [ -z "$XDG_RUNTIME_DIR" ]; then
    export XDG_RUNTIME_DIR="/run/user/$(id -u)"
fi

# Parar e desabilitar o serviço do systemd para o usuário
systemctl --user stop md-tracker-search-provider.service || true
systemctl --user disable md-tracker-search-provider.service || true
systemctl --user daemon-reload

# Remover arquivos relacionados ao search provider
rm -f ~/.local/share/gnome-shell/search-providers/md-tracker-provider.ini
rm -f ~/.local/share/applications/md-tracker-search-provider.desktop
