# DEBIAN/prerm
#!/bin/bash
set -e

# Obter o nome do usuário do diretório HOME padrão
USER=$(logname)
if [ -z "$USER" ]; then
    echo "Não foi possível determinar o usuário atual."
    exit 1
fi

# Configurar o ambiente do usuário para o systemctl --user
XDG_RUNTIME_DIR="/run/user/$(id -u "$USER")"
export XDG_RUNTIME_DIR

# Parar e desabilitar o serviço no contexto do usuário
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user stop md-tracker-search-provider.service || true"
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user disable md-tracker-search-provider.service || true"
su - "$USER" -c "XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR systemctl --user daemon-reload"

echo "MD Tracker Search Provider foi removido com sucesso."
